# command.cron
PATH=/usr/local/bin:/usr/bin:/bin
* * * * * cd /app && python manage.py check >> /tmp/log/django_commands.log 2>&1
* * * * * cd /app && python manage.py my_custom_command >> /tmp/log/django_commands.log 2>&1
