# command.cron
PATH=/usr/local/bin:/usr/bin:/bin
* * * * * cd /app && newrelic-admin run-program python >> /tmp/log/django_commands.log 2>&1
